# deploy/templates/etc/supervisor/conf.d/cache_invalidator.conf
[program:cache_invalidator]
command={{ keystone_venv_bin | dirname }}/bin/capstone-cache-invalidator
user=nobody
autostart=true
autorestart=true
